# curator4-2030330439.conf
[program:curator4-2030330439]
command=python3 /home/curator4/curator4-2030330439/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-2030330439/test.err.log
stdout_logfile=/home/curator4/curator4-2030330439/test.out.log